# misc/events.supervisor.ini
[program:events]
process_name = events
command = /sites/archive-backend/archive-backend events
directory = /sites/archive-backend/
user = archive
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/archive-backend/logs/events.log
