# misc/links.supervisor.ini
[program:links]
process_name = links
command = /sites/archive-backend/archive-backend links
directory = /sites/archive-backend/
user = archive
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/archive-backend/logs/links.log
