# misc/mdb2es.supervisor.ini
[program:mdb2es_esplorer]
process_name = mdb2es_esplorer
command = /sites/mdb2es/mdb2es esplorer
directory = /sites/mdb2es/
user = archive
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/mdb2es/logs/esplorer.log