# misc/supervisor.ini
[program:archive]
process_name = archive
command = /sites/archive-backend/archive-backend server
directory = /sites/archive-backend/
user = archive
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/archive-backend/logs/archive.log